# deploy/systemd/relaymd-orchestrator.service
[Unit]
Description=RelayMD Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Edit this path to your relaymd repo root.
WorkingDirectory=/home/USER/relaymd
ExecStart=uv run relaymd-orchestrator
Restart=on-failure
Environment=RELAYMD_CONFIG=%h/.config/relaymd/config.yaml
# Environment=RELAYMD_API_TOKEN=change-me
# Environment=INFISICAL_TOKEN=<client_id>:<client_secret>

[Install]
WantedBy=default.target
